# assets/autostart/papirusred.desktop
[Desktop Entry]
Encoding=UTF-8
Version=1.0
Type=Application
Terminal=false
Exec=papirus-folders -C red --theme Papirus-Dark
Name=papirusred
Icon=